# curator4-1663155686.conf
[program:curator4-1663155686]
command=python3 /home/curator4/curator4-1663155686/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1663155686/test.err.log
stdout_logfile=/home/curator4/curator4-1663155686/test.out.log